# setup/conf/supervisor/prod_mqtt-im-mqtt2.conf
[program:mqtt-im]
user=ubuntu
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/uwsgi --ini /home/ubuntu/mqtt-im/scripts/uwsgi/prod/mqtt-im.ini
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/uwsgi-webapp.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/uwsgi-webapp.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:exhook9200]
user=ubuntu
environment=DJANGO_SETTINGS_MODULE=server_api.settings.exhook_prod
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/python /home/ubuntu/mqtt-im/exhook_server.py 9200
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/exhook_9200.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/exhook_9200.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:exhook9300]
user=ubuntu
environment=DJANGO_SETTINGS_MODULE=server_api.settings.exhook_prod
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/python /home/ubuntu/mqtt-im/exhook_server.py 9300
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/exhook_9300.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/exhook_9300.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:exhook9400]
user=ubuntu
environment=DJANGO_SETTINGS_MODULE=server_api.settings.exhook_prod
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/python /home/ubuntu/mqtt-im/exhook_server.py 9400
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/exhook_9400.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/exhook_9400.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:exhook9500]
user=ubuntu
environment=DJANGO_SETTINGS_MODULE=server_api.settings.exhook_prod
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/python /home/ubuntu/mqtt-im/exhook_server.py 9500
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/exhook_9500.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/exhook_9500.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:mqtt-client]
user=ubuntu
environment=DJANGO_SETTINGS_MODULE=server_api.settings.client_prod
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/python /home/ubuntu/mqtt-im/mqtt-client.py
directory=/home/ubuntu/mqtt-im/
stdout_logfile=/var/log/mqtt-im/client.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/client.err        ; stderr log path, NONE for none; default AUTO
killasgroup=true

[program:mqtt-celery-task1]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api worker -l info -Ofair -P eventlet -c 100 -n worker1 -Q celery
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-task1.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-task1.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

[program:mqtt-celery-task2]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api worker -l info -Ofair -P eventlet -c 100 -n worker2 -Q celery
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-task2.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-task2.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

[program:mqtt-celery-task3]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api worker -l info -Ofair -P eventlet -c 100 -n worker3 -Q celery
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-task3.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-task3.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

[program:mqtt-celery-task4]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api worker -l info -Ofair -P eventlet -c 100 -n worker4 -Q celery
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-task4.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-task4.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

[program:mqtt-celery-task5]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api worker -l info -Ofair -P eventlet -c 100 -n worker5 -Q celery
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-task5.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-task5.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

[program:mqtt-celery-beat]
command=/home/ubuntu/.virtualenvs/mqtt-im/bin/celery -A server_api beat
environment=DJANGO_SETTINGS_MODULE="server_api.settings.celery_prod"       ; process environment additions (def no adds)
directory=/home/ubuntu/mqtt-im/            ; directory to cwd to before exec (def no cwd)
stdout_logfile=/var/log/mqtt-im/celery-beat.log        ; stdout log path, NONE for none; default AUTO
stderr_logfile=/var/log/mqtt-im/celery-beat.err        ; stderr log path, NONE for none; default AUTO
user=ubuntu

